# runner-setup/prune-texlive-images.timer
[Unit]
Description=Prune Docker

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
