# runner-setup/prune-texlive-images.service
[Unit]
Description=Prune Docker

[Service]
Type=oneshot
ExecStart=/usr/bin/prune-texlive-images.sh

WorkingDirectory=/tmp
